# templates/com.token-economy.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.token-economy.weekly</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>__HOME__/.claude/scripts/token_economy_weekly.py</string>
        <string>--days</string>
        <string>7</string>
    </array>
    <!-- Every Monday 09:00 (launchd catches up if the machine was asleep) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/token-economy-weekly.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/token-economy-weekly.err</string>
</dict>
</plist>
